# build-r64-arch-utils-git.install
post_upgrade() {
  sudo systemctl --force --no-pager disable setmac.service
  echo
  sudo systemctl --force --no-pager enable setmac.service
  echo
}

post_install() {
  post_upgrade
}

pre_remove() {
  sudo systemctl --force --no-pager disable setmac.service
  rm -vf  /boot/dtbos/*-fixed.dtb
}
